# docs/jarvis_delta_ingest_v2/systemd/jarvis-delta-cold.timer
[Unit]
Description=Jarvis Delta Ingest (cold) daily at 03:00

[Timer]
OnCalendar=*-*-* 03:00:00
Unit=jarvis-delta@cold.service

[Install]
WantedBy=timers.target
