# docs/jarvis_delta_ingest_v2/systemd/jarvis-delta-warm.timer
[Unit]
Description=Jarvis Delta Ingest (warm) every 4 hours

[Timer]
OnBootSec=5m
OnUnitActiveSec=4h
Unit=jarvis-delta@warm.service

[Install]
WantedBy=timers.target
